# resources/acd.service
[Unit]
Description=Agent Console Dashboard Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/acd daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
